# debian/postinst
#!/bin/sh
# postinst script for kgp-database
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

#
# Divert file, originally by azatoth in bubba-backend package
#
# Usage:
#
# divert package orig_file orig_file_diverted new_source_file service
#
# Where
#	package   		- package doing diversion, its files wont be diverted
#	orig_file 		- Which file should be diverted
#	orig_file_diverted	- Where to should file be diverted/saved
#	new_sourde_file		- Which file should be used as replacement
#	service			- service to restart after diversion or "false"
#
# Sideeffects, if function divert_cb_"service" function exists it will be called
# instead of service restart being performed
#
divert() {
	package=$1;
	file=$2;
	diverted=$3;
	source_file=$4;
	service=$5;
	if [ x$(dpkg-divert --truename $file) = x$file ]; then # we are not diverted at all
		dpkg-divert --add --package $package --rename --divert $diverted $file;
		if [ ! -e $file ]; then
			if [ "false" != "$source_file" ]; then
				ln -sf $source_file $file;
			fi
			if $(command -v "divert_cb_$service" >/dev/null); then 
				divert_cb_$service "$@";
			elif [ "$service" != "false" ]; then
				systemctl --quiet try-restart $service || true;
			fi
		fi
	fi

}


case "$1" in
	configure)
		#Disable services that should not run by default
		for service in mariadb mysql
		do
			if systemctl --quiet is-enabled $service
			then
				systemctl disable $service
				systemctl stop $service
			fi
		done
	

		# Fixups
		if ! [ -d /var/opi/mysql ]
		then
			mkdir -p /var/opi/mysql
			chown -h -R mysql:mysql /var/opi/mysql

			# Initialise db-storage on new location
			/usr/bin/mysql_install_db --skip-test-db --auth-root-authentication-method=socket --rpm --cross-bootstrap --user=mysql --disable-log-bin 2>&1
		fi


    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
